# config/home/service/clash/run
#!/bin/bash

/root/go/bin/clash_pro -f /root/.config/clash/tw.yaml
