# service/caddy/run
#!/bin/bash

exec caddy file-server --browse --root ~/d/ 2>&1
