# config/etc/sv/runsvdir-user/run
#!/bin/sh

export USER="root"
export HOME="/root"

groups="$(id -Gn "$USER" | tr ' ' ':')"
svdir="$HOME/service"

exec chpst -u "$USER:$groups" runsvdir "$svdir"
