# service/clash/run
#!/bin/bash

# wget -O /tmp/jms.conf 'http://sub.totorow.xyz:25500/getprofile?name=profiles/jms.ini&token=xxx'
exec ~/go/bin/clash -f /tmp/jms.conf
